# Slides/Unit-Files/Path/org.cups.cupsd.path
[Unit]
Description=CUPS Scheduler
PartOf=org.cups.cupsd.service

[Path]
PathExistsGlob=/var/spool/cups/d*

[Install]
WantedBy=multi-user.target
