# Documentation/Examples/Automount/mnt-resources.automount
# mnt-resources.automount defines an automount point
# at /mnt/resources folder
# Bruno Mondelo Giaramita
# 2017-05-12 Escola Del Treball

[Unit]
Description=Resource folder automount point
Documentation=https://github.com/mondelob

[Automount]
Where=/mnt/resources
DirectoryMode=750
TimeoutIdleSec=2min 30s
